# ansible/roles/tls/templates/refresh-certs.timer.j2
[Unit]
Description=Refresh certs for {{ tls_domain }} daily

[Timer]
OnCalendar=daily
Persistent=true
# 6 hour randomized delay
RandomizedDelaySec=21600

[Install]
WantedBy=timers.target
